# configs/systemd/vdr.service
[Unit]
Description=VDR
After=lirc.service
Wants=lirc.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/rc.d/vdr start
ExecStop=/etc/rc.d/vdr stop
KillMode=none

[Install]
WantedBy=multi-user.target
